# backup.rules
#backup platte
BUS=="usb", KERNEL=="sd?1", SYSFS{serial}=="2011061900872", SYMLINK+="backup", RUN+="/usr/local/sbin/backup.sh"
